# mind/etc/local.d/00-ntp-client.start
#!/bin/bash

(while [ 1 ] ; do
    /etc/init.d/ntp-client restart
    sleep 1m
done) &
