# mind/etc/local.d/00-ntp-client.stop
#!/bin/bash

kill $(ps uax | grep 00-ntp-client.start | grep -v grep | awk '{print $2}')
