# mind/etc/local.d/05-co2mond.stop
#!/bin/sh

DATA_DIR="/tmp/co2mond"
killall co2mond 2>>/dev/null
rm -rf "${DATA_DIR}"
