# desktop/etc/local.d/01-fix-wifi.start
#!/bin/bash

logger "running $0"

(while [ 1 ] ; do
    rfkill unblock 0
    rfkill unblock 1
    rfkill block bluetooth
    sleep 3s
done) &
