# desktop/etc/local.d/03-cpu-temp-control.start
#!/bin/bash

logger "running $0"

MAX_ALLOWED_TEMP=90

(while [ 1 ] ; do
    MAX_CPU_TEMP=$(sensors | grep '^Core' | awk '{print $3}' | sed 's/\.[0-9]//g;s/[^0-9]*//g' | sort | tail -n1)
    if [ "${MAX_CPU_TEMP}" -ge "${MAX_ALLOWED_TEMP}" ]; then
        /usr/local/sbin/switch-governor.sh powersave
        sleep 1m
    else
        /usr/local/sbin/switch-governor.sh ondemand
    fi
    sleep 3s
done) &
