# mlaptop/etc/local.d/00-fix-ssd-issues.stop
#!/bin/bash

kill $(ps uax | grep 00-fix-ssd-issues.start | grep -v grep | awk '{print $2}')
