# mlaptop/etc/local.d/02-keyboard.start
#!/bin/bash

# http://askubuntu.com/questions/370857/cant-adjust-screen-brightness-on-macbook-pro-10-1-ubuntu-13-10 
setpci -v -H1 -s 00:01.00 BRIDGE_CONTROL=0

chmod a+w /sys/devices/platform/applesmc.768/leds/smc\:\:kbd_backlight/brightness
chmod a+w /sys/class/backlight/acpi_video0/brightness
