# mlaptop/etc/local.d/03-cpu-gov-control.stop
#!/bin/bash

kill $(ps uax | grep 03-cpu-gov-control.start | grep -v grep | awk '{print $2}')
