# lenlap/etc/udev/rules.d/95-uvc-webcams.rules
ACTION=="add", ATTRS{idVendor}=="1532", RUN+="/usr/local/sbin/load-webcam-settings.sh"
