# common/etc/udev/rules.d/99-keyboard.rules
KERNEL=="event[0-9]*", ACTION=="add", RUN+="/usr/local/sbin/set-keyboard-layout.sh"
