# lenlap/etc/local.d/90-watch-compiler-processes.stop
#!/bin/bash

kill $(ps uax | grep 90-watch-compiler-processes.start | grep -v grep | awk '{print $2}')
