# lenlap/etc/local.d/99-mpd-broken-pipe-workaround.stop
#!/bin/bash

kill $(ps uax | grep 99-mpd-broken-pipe-workaround.start | grep -v grep | awk '{print $2}')
