# lenlap/etc/local.d/99-mpd-broken-pipe-workaround.start
#!/bin/bash

MPD_LOG=$(ls --color=never -t1 /var/lib/mpd/log-* | grep -v gz$ | sort -r)

logger "running $0"

(
    tail -n1 -F "${MPD_LOG}" | while read line ; do echo "${line}" | grep 'Broken pipe' >> /dev/null && sleep 0.2s && mpc play ; done
) &
