# lenlap/etc/local.d/99-mem-watcher.stop
#!/bin/bash

kill $(ps uax | grep 99-mem-watcher.start | grep -v grep | awk '{print $2}')
