# lenlap/etc/local.d/99-non-blinking-block-cursor.start
#!/bin/bash

# non-blinking block cursor
for i in {1..12}
do
    echo -ne '\033[?17;0;121c' > "/dev/tty${i}"
done
