# lenlap/etc/local.d/local.stop
#!/bin/sh

kill $(ps uax | grep discharge_watcher.py | grep -v grep | awk '{print $2}')
