# etc/local.d/01-fix-wifi.start
#!/bin/bash

rfkill block 0
rfkill block bluetooth
rfkill unblock 1
